# prometheus/roles/prometheus/templates/prometheus.service
[Unit]
Description=prometheus

[Service]
ExecStart=/opt/prometheus-{{prometheus_version}}.linux-amd64

[Install]
WantedBy=multi-user.target